# puppetmaster/supervisord.conf
[supervisord]
nodaemon=true

[program:puppet]
command=/usr/bin/puppet master --no-daemonize --verbose
